# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true

[program:phpfpm]
command=php-fpm
autostart=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true